# backend/deploy/dpf-mqtt-ingestor.service.example
[Unit]
Description=DPF tracker MQTT ingestor
After=network-online.target postgresql.service mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Documents/Ford-Mondeo-MK-DPF-Tracker
EnvironmentFile=/etc/dpf-backend.env
ExecStart=/home/pi/Documents/Ford-Mondeo-MK-DPF-Tracker/backend/.venv/bin/dpf-mqtt-ingestor --log-level INFO
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
